# jetauto_bringup/service/jupyter.service
[Unit]
Description=Jupyter Notebook

[Service]
Type=simple
User=jetauto
ExecStart=/usr/local/bin/jupyter-lab --allow-root --port 8888 --no-browser
WorkingDirectory=/home/jetauto/
StandardOutput=null
StandardError=null

[Install]
WantedBy=default.target
